# docker/service/dashing/run
#!/bin/sh

cd /app
exec setuser app dashing start
